# mac/LaunchAgents/fastri-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>fastri-server</string>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/local/bin/ruby</string>
      <string>/opt/local/bin/fastri-server</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/takeru/srv/fastri-server</string>

    <key>KeepAlive</key>
    <true/>

    <key>StandardErrorPath</key>
    <string>/Users/takeru/var/log/fastri-server.log</string>

    <key>StandardOutPath</key>
    <string>/Users/takeru/var/log/fastri-server.log</string>

  </dict>
</plist>
